# deploy/openwrt/files/etc/init.d/mjolnir-babeld
#!/bin/sh /etc/rc.common
# procd service for babeld, driven by mjolnir-meshd's rendered config.
# Split of concerns (mjolnir-mesh-m8t / qz9): meshd RENDERS the config
# (/etc/mjolnir/babeld.conf); procd OWNS the process AND restarts babeld when
# that file changes (procd_set_param file, below). meshd no longer drives the
# restart itself — that synchronous procd-restart loop wedged meshd (qz9). meshd
# never fork()s babeld itself (that chain orphaned babelds on SIGKILL). Pairs with
# /etc/init.d/mjolnir-meshd. Install: enable this, and DISABLE the stock babeld
# service so the two don't both run.

USE_PROCD=1
START=96   # after mjolnir-meshd (95) has rendered the conf + brought up the L2
STOP=09

PROG=/usr/sbin/babeld
CONF=/etc/mjolnir/babeld.conf

# Kill any babeld on OUR config that procd does not own (mjolnir-mesh-nrr).
# Interrupted/raced restarts have left a TERM-surviving babeld behind the new
# procd instance — two babelds fight over kernel routes and the loser's stale
# routes linger. Runs on every start: whatever race created the stray, the
# next (re)start reaps it. procd-owned pids are excluded via ubus, so a plain
# `start` on an already-running service never kills the healthy instance.
#
# Match by EXACT full cmdline via /proc, never pgrep -f: substring matching
# kills innocent processes whose command text merely CONTAINS the babeld path
# (proven live: it reaped the operator's own SSH diagnostic shell).
reap_stray_babelds() {
	local owned p pid cmd want
	owned=$(ubus call service list '{"name":"mjolnir-babeld"}' 2>/dev/null \
		| jsonfilter -e '@["mjolnir-babeld"].instances[*].pid' 2>/dev/null)
	want="$PROG -c $CONF"
	for p in /proc/[0-9]*; do
		cmd=$(tr '\0' ' ' < "$p/cmdline" 2>/dev/null)
		[ "${cmd% }" = "$want" ] || continue
		pid="${p#/proc/}"
		case " $owned " in
			*" $pid "*) ;;
			*)
				logger -t mjolnir-babeld "reaping stray babeld pid $pid (not procd-owned)"
				kill "$pid" 2>/dev/null
				;;
		esac
	done
}

# True (0) if a babeld running OUR exact config is currently alive. Same
# exact-cmdline /proc match as reap_stray_babelds — never pgrep -f, whose
# substring matching reaps innocent processes (see the reap note above).
babeld_running() {
	local p cmd want
	want="$PROG -c $CONF"
	for p in /proc/[0-9]*; do
		cmd=$(tr '\0' ' ' < "$p/cmdline" 2>/dev/null)
		[ "${cmd% }" = "$want" ] && return 0
	done
	return 1
}

# Flush routes a dead babeld left behind (mjolnir-mesh-bii). babeld installs
# mesh routes with proto 'babel'; on a crash / SIGKILL / interrupted restart
# those routes LINGER in the kernel FIB — reap_stray_babelds reaps the stray
# PROCESS but not its routes. A lingering 'proto babel' default is then
# (a) re-exported by the next babeld as a foreign kernel route, breaking chj's
# split-horizon invariant so a no-WAN node advertises a bogus 0.0.0.0/0, and
# (b) learnable by the real gateway OVER its own DHCP WAN default, blackholing
# egress back into the mesh (mjolnir-mesh-5lw). Before this, only a reboot
# cleared it. Safe: the fresh babeld re-installs every route it still believes
# in within one update cycle. GUARDED on babeld_running so we never yank routes
# out from under a healthy instance (e.g. a no-op `start` while already up) —
# only flush when we are genuinely bringing babeld up with no live instance.
flush_stale_babel_routes() {
	babeld_running && return 0
	logger -t mjolnir-babeld "flushing stale proto-babel routes before (re)start"
	ip route flush proto babel 2>/dev/null
	ip -6 route flush proto babel 2>/dev/null
}

start_service() {
	# meshd renders CONF and (re)starts us once it has an interface to route over.
	# If it isn't there yet (early boot), do nothing — meshd will start us when it is.
	[ -f "$CONF" ] || return 0
	reap_stray_babelds
	flush_stale_babel_routes
	procd_open_instance
	procd_set_param command "$PROG" -c "$CONF"
	# procd watches the meshd-rendered config and restarts babeld on change
	# (mjolnir-mesh-qz9) — so meshd never has to drive the restart synchronously.
	procd_set_param file "$CONF"
	procd_set_param respawn 3600 5 0
	procd_set_param stdout 1
	procd_set_param stderr 1
	procd_close_instance
}
