# deploy/openwrt/files/etc/init.d/mjolnir-babeld
#!/bin/sh /etc/rc.common
# procd service for babeld, driven by mjolnir-meshd's rendered config.
# Split of concerns (mjolnir-mesh-m8t / qz9): meshd RENDERS the config
# (/etc/mjolnir/babeld.conf); procd OWNS the process AND restarts babeld when
# that file changes (procd_set_param file, below). meshd no longer drives the
# restart itself — that synchronous procd-restart loop wedged meshd (qz9). meshd
# never fork()s babeld itself (that chain orphaned babelds on SIGKILL). Pairs with
# /etc/init.d/mjolnir-meshd. Install: enable this, and DISABLE the stock babeld
# service so the two don't both run.

USE_PROCD=1
START=96   # after mjolnir-meshd (95) has rendered the conf + brought up the L2
STOP=09

PROG=/usr/sbin/babeld
CONF=/etc/mjolnir/babeld.conf

# Kill any babeld on OUR config that procd does not own (mjolnir-mesh-nrr).
# Interrupted/raced restarts have left a TERM-surviving babeld behind the new
# procd instance — two babelds fight over kernel routes and the loser's stale
# routes linger. Runs on every start: whatever race created the stray, the
# next (re)start reaps it. procd-owned pids are excluded via ubus, so a plain
# `start` on an already-running service never kills the healthy instance.
#
# Match by EXACT full cmdline via /proc, never pgrep -f: substring matching
# kills innocent processes whose command text merely CONTAINS the babeld path
# (proven live: it reaped the operator's own SSH diagnostic shell).
reap_stray_babelds() {
	local owned p pid cmd want
	owned=$(ubus call service list '{"name":"mjolnir-babeld"}' 2>/dev/null \
		| jsonfilter -e '@["mjolnir-babeld"].instances[*].pid' 2>/dev/null)
	want="$PROG -c $CONF"
	for p in /proc/[0-9]*; do
		cmd=$(tr '\0' ' ' < "$p/cmdline" 2>/dev/null)
		[ "${cmd% }" = "$want" ] || continue
		pid="${p#/proc/}"
		case " $owned " in
			*" $pid "*) ;;
			*)
				logger -t mjolnir-babeld "reaping stray babeld pid $pid (not procd-owned)"
				kill "$pid" 2>/dev/null
				;;
		esac
	done
}

# True (0) if a babeld running OUR exact config is currently alive. Same
# exact-cmdline /proc match as reap_stray_babelds — never pgrep -f, whose
# substring matching reaps innocent processes (see the reap note above).
babeld_running() {
	local p cmd want
	want="$PROG -c $CONF"
	for p in /proc/[0-9]*; do
		cmd=$(tr '\0' ' ' < "$p/cmdline" 2>/dev/null)
		[ "${cmd% }" = "$want" ] && return 0
	done
	return 1
}

# Clear stale mesh-learned DEFAULT routes before (re)starting babeld
# (mjolnir-mesh-bii). babeld installs the learned internet default via a MESH
# interface (br-mesh, or mjolnir0 in overlay mode); on a crash / SIGKILL /
# interrupted restart it LINGERS in the kernel FIB — reap_stray_babelds reaps the
# stray PROCESS but not its routes. A lingering mesh default is then
# (a) re-exportable by the next babeld, breaking chj's split-horizon invariant so
# a no-WAN node advertises a bogus 0.0.0.0/0, and (b) learnable by the real
# gateway OVER its own WAN default, blackholing egress back into the mesh
# (mjolnir-mesh-5lw). Before this, only a reboot cleared it.
#
# We must NOT filter by `proto babel`: OpenWrt's busybox `ip` SILENTLY IGNORES the
# proto selector (verified on the mt7981 fleet — `ip route show proto babel`
# returns the whole table), so `ip route flush proto babel` would flush EVERY
# route and cut the node off. Match surgically by output interface instead: only
# a mesh-interface DEFAULT is ever babel-learned; a real uplink default is out
# eth0/wan and link routes are not defaults, so neither is ever touched. Safe: the
# fresh babeld re-installs a still-valid learned default within one update cycle.
# GUARDED on babeld_running so a healthy instance (e.g. a no-op `start`) is never
# disturbed — only clean up when bringing babeld up with no live instance.
flush_stale_babel_routes() {
	babeld_running && return 0
	# Lines look like: default via 10.254.12.214 dev br-mesh onlink
	ip -4 route show default 2>/dev/null | while read -r kw _via gw devkw ifname _rest; do
		[ "$kw" = default ] && [ "$devkw" = dev ] || continue
		case "$ifname" in
			br-mesh|mjolnir0)
				logger -t mjolnir-babeld "clearing stale mesh default before restart: via $gw dev $ifname"
				ip -4 route del default via "$gw" dev "$ifname" 2>/dev/null
				;;
		esac
	done
}

start_service() {
	# meshd renders CONF and (re)starts us once it has an interface to route over.
	# If it isn't there yet (early boot), do nothing — meshd will start us when it is.
	[ -f "$CONF" ] || return 0
	reap_stray_babelds
	flush_stale_babel_routes
	procd_open_instance
	procd_set_param command "$PROG" -c "$CONF"
	# procd watches the meshd-rendered config and restarts babeld on change
	# (mjolnir-mesh-qz9) — so meshd never has to drive the restart synchronously.
	procd_set_param file "$CONF"
	procd_set_param respawn 3600 5 0
	procd_set_param stdout 1
	procd_set_param stderr 1
	procd_close_instance
}
